# etc/supervisor.ini
[program:frontend]
environment=DAV_BASE=/home/afsdav
command=/home/afsdav/sbin/lighttpd -D -f /home/afsdav/etc/frontend.conf -m /home/afsdav/lib/lighttpd/

[program:angel]
environment=DAV_BASE=/home/afsdav
command=/home/afsdav/sbin/angel
umask=000
user=nobody
